# rites/pamusb/templates/udev-lock.rules
# pamusb: trigger session lock/unlock on USB add/remove
# Generated by pamusb setup — do not edit manually
# Watches the configured main (non-LUKS) partition identity for plug/unplug events

ACTION=="add|remove", SUBSYSTEM=="block", ENV{MATCH_ENV}=="{{MATCH_VALUE}}", RUN+="/usr/local/bin/usb-pam-guard.sh %E{ACTION}"
