# spells/usb-auth/templates/udev-lock.rules
# usb-auth: trigger session lock/unlock on USB add/remove
# Generated by usb-auth setup — do not edit manually
# Watches the main (non-LUKS) partition UUID for plug/unplug events

ACTION=="add|remove", SUBSYSTEM=="block", ENV{ID_FS_UUID}=="{{MAIN_UUID}}", RUN+="/usr/local/bin/usb-pam-guard.sh %E{ACTION}"
